# templates/airflow-webserver.service.erb
[Unit]
Description= Airflow webserver
Wants=basic.target
After=basic.target network.target

[Service]
PIDFile = <%= scope.lookupvar('run_folder') %>/airflow_webserver.pid
ExecStart=/usr/local/bin/gunicorn -w <%= scope.lookupvar('gunicorn_workers') %> -t 120 -b <%= scope.lookupvar('web_server_host') %>:<%= scope.lookupvar('web_server_port') %> airflow.www.app:cached_app()
ExecReload=/bin/kill -HUP $MAINPID
Environment="AIRFLOW_HOME=<%= scope.lookupvar('home_folder') %>"
KillMode=process
Restart=on-failure
RestartSec=15s
User=<%= scope.lookupvar('user') %>
Group=<%= scope.lookupvar('group') %>

[Install]
WantedBy=multi-user.target





